# ansible/roles/docker-install/templates/docker-prune.service.j2
[Unit]
Description=Prune unused Docker images

[Service]
Type=simple
ExecStart=/usr/bin/docker system prune --volumes --force

[Install]
WantedBy=default.target
